# poller/com.japan-trip.poller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.japan-trip.poller</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/jakeroyston/documents/github/japan-trip-ios/poller/run.sh</string>
  </array>

  <!-- Run every 30 minutes. Unlike cron, launchd queues missed runs
       and fires them on wake, so the Mac mini sleeping is not a problem. -->
  <key>StartInterval</key>
  <integer>1800</integer>

  <!-- Run once immediately when the agent is loaded -->
  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/jakeroyston/documents/github/japan-trip-ios/poller/poller.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/jakeroyston/documents/github/japan-trip-ios/poller/poller.log</string>

  <!-- Restart if the script crashes -->
  <key>KeepAlive</key>
  <false/>
</dict>
</plist>
